# roles/node_exporter/templates/node_exporter.service.j2
############################################
#                                          #
#            Managed by Ansible            #
#           Do not edit manually           #
#                                          #
############################################

[Unit]
Description=NodeExporter

[Service]
TimeoutStartSec=0
User=node_exporter
ExecStart={{ node_exporter_bin_path }} --web.listen-address={{ node_exporter_host }}:{{ node_exporter_port }} {% if node_exporter_textfile_collector_dir is defined %}--collector.textfile.directory={{ node_exporter_textfile_collector_dir }}{% endif %}{{ node_exporter_options }}

[Install]
WantedBy=multi-user.target