# roles/prometheus_smart_exporter/templates/prometheus_smart_helper.service.j2
############################################
#                                          #
#            Managed by Ansible            #
#           Do not edit manually           #
#                                          #
############################################

[Unit]
Description=Helper service to provide S.M.A.R.T. info to prometheus_smart_exporter
Requires=prometheus_smart_helper.socket

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/smart_exporter_helper -vv {% if prometheus_smart_exporter_ignore is defined %}--ignore {{ prometheus_smart_exporter_ignore }} {% endif %}

PrivateTmp=true
PrivateDevices=no
ProtectHome=yes
PrivateNetwork=yes